# deploy/lcg.service
# systemd unit: on every boot, point DNS at this instance's new public IP, then
# bring up all three LCG app stacks. Because EBS persists across stop/start, the
# compiled Elixir _build, Postgres data, and the MySQL/MariaDB volumes survive,
# so this is a fast container start rather than a full rebuild.
#
# Apps started:
#   - DragnCards: /srv/dragncards/compose.prod.yml  (postgres + backend + Caddy)
#   - RingsDB:    /srv/ringsdb/docker-compose.yml    (mariadb + php/apache :8001)
#   - MarvelCDB:  /srv/marvelsdb/docker-compose.yml  (mysql + php/apache :8000)
#   - ashteki:    /srv/ashteki/compose.prod.yml      (mongo + redis + lobby :4100 + gamenode :10000)
# Caddy (in the DragnCards stack, host networking) reverse-proxies the *DB and
# ashes subdomains to the host ports those containers publish on loopback.
#
# Install:
#   sudo cp deploy/lcg.service /etc/systemd/system/lcg.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now lcg.service
#   # (if migrating from the old unit: sudo systemctl disable --now dragncards.service)
[Unit]
Description=LCG apps (Route53 update + DragnCards/RingsDB/MarvelCDB/ashteki docker compose stacks)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/srv/dragncards
# 1) Update Route 53 A records (lcg/ringsdb/marvelcdb) to this instance's current IP.
ExecStartPre=/srv/dragncards/deploy/update-route53.sh
# 2) Start the DragnCards prod stack (postgres + backend + Caddy).
ExecStart=/usr/bin/docker compose -f /srv/dragncards/compose.prod.yml up -d
# 3) Start the RingsDB stack.
ExecStart=/usr/bin/docker compose -f /srv/ringsdb/docker-compose.yml up -d
# 4) Start the MarvelCDB stack.
ExecStart=/usr/bin/docker compose -f /srv/marvelsdb/docker-compose.yml up -d
# 5) Start the ashteki stack (mongo + redis + lobby + gamenode).
ExecStart=/usr/bin/docker compose -f /srv/ashteki/compose.prod.yml up -d
# Stop all stacks cleanly on shutdown (reverse order).
ExecStop=/usr/bin/docker compose -f /srv/ashteki/compose.prod.yml down
ExecStop=/usr/bin/docker compose -f /srv/marvelsdb/docker-compose.yml down
ExecStop=/usr/bin/docker compose -f /srv/ringsdb/docker-compose.yml down
ExecStop=/usr/bin/docker compose -f /srv/dragncards/compose.prod.yml down

[Install]
WantedBy=multi-user.target
